# iso-profiles-orig/community/sway/desktop-overlay/etc/skel/.config/autostart/fcitx5.desktop
[Desktop Entry]
Exec=/usr/bin/fcitx5 -r
Name=fcitx5
Comment=Input Method
Type=Application
