# iso-profiles-orig/manjaro/gnome/desktop-overlay/etc/skel/.config/autostart/manjaro-hello.desktop
Type=Application
Categories=GNOME;GTK;System;
StartupNotify=false
Name=Manjaro Hello
Exec=/usr/bin/manjaro-hello
Icon=manjaro
X-GNOME-Autostart-Delay=8
